# default/app/nodejs/supervisord.conf
[program:node]
command = /home/dotcloud/code/node_modules/sm/bin/sm-cli run --mode production
directory = /home/dotcloud/current